# rootfs/etc/services/hostapd/run
#!/bin/sh
exec hostapd /tmp/hostapd.conf
