# rootfs/etc/services/dhcpd/run
#!/bin/sh
exec udhcpd -f /etc/udhcpd.conf
